# web/deploy/backup/kyber-web-backup.service
[Unit]
Description=Kyber-Zip web encrypted backup
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kyber-web-backup
